# crontab
30 * * * * sh /home/pi/garden/current/take-current-snapshot 2>&1 /home/pi/snapshots.log